# tmcd/freebsd6/netif-emulab.in
#
# Emulab version of netif to DHCP on all interfaces simultaneously.
# Runs before the standard netif, and effectively neuters that script.
# Note the dependencies below.  Every script that declares that it must run
# BEFORE netif is listed as a REQUIRE here so that we are properly ordered
# w.r.t. those scripts.
#

# PROVIDE: netif-emulab
# REQUIRE: atm1 ipfilter mountcritlocal pccard serial sppp sysctl
# REQUIRE: adjkerntz hostname ipfs ipnat kldxref random
# BEFORE: netif
# KEYWORD: nojail

. /etc/rc.subr
. /etc/network.subr

name="netif-emulab"
start_cmd="cnet_start"
stop_cmd="cnet_stop"
_cmdifn=

# some emulab paths
ELAB_LOGDIR=@CLIENT_VARDIR@/logs
ELAB_BOOTDIR=@CLIENT_VARDIR@/boot

# this is a separate function so we can redirect all the output below
cnet_dhcp()
{
	#
	# Argh! FreeBSD >= 6 dhclient also allows only a single interface
	# to be specified.  I tried to make this work by launching multiple,
	# simultaneous dhclients, but that caused EEPROM read failures on
	# em devices.  So we will use the "pure" port version if it is
	# installed.
	#
	if [ -x /usr/local/sbin/dhclient ]; then
	    echo "Using dhclient port..."
	    dhclient_program="/usr/local/sbin/dhclient"
	    ${dhclient_program} ${dhclient_flags} $*
	else
	    echo "Using default dhclient..."

	    #
	    # Argh redux! FreeBSD >= 6 has rewritten dhclient and it no longer
	    # takes the -i option.  So we remove that option here.
	    #
	    dhclient_flags=`echo $dhclient_flags | sed -e 's/-i [0-9][0-9]*//'`

	    for _if in $*; do
		${dhclient_program} ${dhclient_flags} -b $_if
	    done
	    wait
	fi
}

cnet_start()
{
	rm -f $ELAB_BOOTDIR/controlif

	#
	# Find all the relevant networks IFs.
	# Don't do loopback, funky tunnel devices
	# (e.g., IP-over-string-and-cans).
	#
	# At Utah Emulab we also exclude wireless interfaces as our
	# wireless machines are in a production ether and might get
	# an answer!
	#
	_ifs=""
	network_interfaces=`ifconfig -l | sed -e 's/lo0//'`
	for _if in $network_interfaces; do
	case $_if in
	lo*|gif*|faith*|tun*|plip*)
		;;
	ath*)
		;;
	*)
		_ifs="$_ifs $_if"
		;;
	esac
	done

	# 
	# For the widearea case: decide if we want static or dynamic 
	# network configuration based on a config file.
	# XXX: should try to look on other media as well, instead of just
	# the boot dongle.
	#
	WACONFIG=/etc/emulab/waconfig
	if [ -e "$WACONFIG" ]; then
	        # key=val pairs become shell vars
		. "$WACONFIG"

		echo "Found Emulab widearea config info..."

		if [ "$WA_BOOTMETHOD" = "static" -a "$WA_MAC" != "" ]; then 
			# try to find an iface matching WA_MAC
			cnetif=""
			for _if in $network_interfaces; do
			        mac=`ifconfig $_if | grep ether | sed -e 's/[[:space:]]*ether[[:space:]]*//'`
				echo "$mac" | grep -i "$WA_MAC"
				if [ "$?" = "0" ]; then
				        cnetif="$_if"
				        break
				fi
			done

			if [ x"$cnetif" = x"" ]; then
			        echo -n "Could not find iface with "
				echo "$WA_MAC; trying DHCP!"
			        WA_BOOTMETHOD="dhcp"
			elif [ x"WA_HOSTNAME" = "x" \
			        -o x"WA_DOMAIN" = "x" \
			        -o x"WA_IP_ADRR" = "x" \
			        -o x"WA_IP_NETMASK" = "x" \
			        -o x"WA_IP_GATEWAY" = "x" \
			        -o x"WA_IP_DNS1" = "x"  ]; then
				echo -n "Missing static IP config "
				echo "vars; trying DHCP!"
				WA_BOOTMETHOD="dhcp"
			fi
		fi
	else
		WA_BOOTMETHOD="dhcp"
	fi

	if [ "$WA_BOOTMETHOD" = "dhcp" ]; then
	        #
	        # use good ol' Emulab dhcp!
		#
		echo "Emulab looking for control net among: $_ifs ..."
		cnet_dhcp $_ifs >$ELAB_LOGDIR/netif-emulab.log 2>&1
		if [ -e $ELAB_BOOTDIR/controlif ]; then
			echo "Emulab control net is `cat $ELAB_BOOTDIR/controlif`"
		else
			echo "*** No Emulab control net found!"
		fi
	else
		#
		# configure according to static WA_X vars:
		#
		echo "Emulab statically configuring control net on $cnetif ..."
		ifconfig "$cnetif" inet "$WA_IP_ADDR" netmask "$WA_IP_NETMASK" up
		route add default "$WA_IP_GATEWAY"

		# bring up lo0 too, since we neuter netif
		ifconfig lo0 inet 127.0.0.1 up
		route add 127.0.0.1 -interface lo0

		# setup resolv.conf
		echo "search $WA_DOMAIN" > /etc/resolv.conf
		echo "nameserver $WA_IP_DNS1" >> /etc/resolv.conf
		if [ x"$WA_IP_DNS2" != x ]; then
			echo "nameserver $WA_IP_DNS2" >> /etc/resolv.conf
		fi

		# set hostname
		hostname "$WA_HOSTNAME"

		# setup hosts file
		echo "$WA_IP_ADDR ${WA_HOSTNAME}.${WA_DOMAIN} $WA_HOSTNAME" >> \
		     /etc/hosts

		# setup a few necessary emulab files...
		echo "$cnetif" > $ELAB_BOOTDIR/controlif
		if [ -e "/etc/emulab/bossnode" ]; then
		        bossnode=`cat /etc/emulab/bossnode`
			echo `host "$bossnode"` > $ELAB_BOOTDIR/bossip
		fi
		echo "$WA_HOSTNAME" > $ELAB_BOOTDIR/realname
		echo "$WA_IP_GATEWAY" > $ELAB_BOOTDIR/routerip
		echo "$WA_IP_ADDR" > $ELAB_BOOTDIR/myip
		echo "$WA_IP_NETMASK" > $ELAB_BOOTDIR/mynetmask
	fi
}

cnet_stop()
{
	rm -f $ELAB_BOOTDIR/controlif
}

load_rc_config $name
run_rc_command $*
