# install/ports/tftp-hpa/pkg-plist
libexec/tftpd
man/man8/tftpd.8
etc/rc.d/tftpd-hpa.sh
etc/tftpd.rules
