# tmcd/freebsd/jail/rc.local
#!/bin/sh
#
# The point of this is to startup testbed stuff inside the jail.
#
. /etc/emulab/paths.sh

if [ -x $BINDIR/jailctl ]; then
    $BINDIR/jailctl start
fi

exit 0

