# clientside/tmcc/freebsd5/frisbee/rc.local
#!/bin/sh

if [ -x /etc/testbed/rc.frisbee ]; then
       	echo "Playing Frisbee ..."
        /etc/testbed/rc.frisbee
fi
