# clientside/tmcc/freebsd/jail/rc.local
#!/bin/sh
#
# The point of this is to startup testbed stuff inside the jail.
#
. /etc/emulab/paths.sh

if [ -x $BINDIR/rc/rc.injail -a ! -e $ETCDIR/nojail ]; then
    $BINDIR/rc/rc.injail
fi

exit 0

